# files/system/etc/profile.d/cloche-fetch.sh
if command -v fastfetch &> /dev/null; then
    fastfetch --config /etc/fastfetch/config.jsonc
fi
